# provisioning/roles/assetsservers/templates/etc/supervisor/conf.d/brunch.conf.j2
[program:brunch]
command = brunch watch
autorestart = true
startretries = 65535
user = {{ user }}
directory = {{ root }}
