# vendor/docker/etc/services.d/fbserve/run
#!/bin/execlineb -P

with-contenv
cd /src
s6-setuidgid fanboi2

/src/.venv/bin/fbctl serve