# provisioning/files/etc/supervisor/conf.d/celery.conf.j2
[program:celery]
command = ${virtualenv}/bin/fb2_celery ${base}/settings.ini worker
autorestart = true
startretries = 65535
user = {{ user }}
directory = {{ root }}
