# provisioning/files/etc/supervisor/conf.d/uwsgi.conf.j2
[program:uwsgi]
command = uwsgi --master --emperor /etc/uwsgi
autorestart = true
startretries = 65535
user = root
